# s6-overlay/s6-rc.d/casaos-main/run
#!/command/with-contenv bash

# Get the message bus URL
runtime_file="/var/run/casaos/message-bus.url"

while [ ! -f "$runtime_file" ]; do
    sleep 1
    echo "Casa Main waiting for message bus URL..."
done

# Wait for /var/run/casaos/routes.json to be created
while [ ! -f /var/run/casaos/routes.json ]; do
    echo "Waiting for /var/run/casaos/routes.json to be created..."
    sleep 1
done

# Define log filter function
filter_logs() {
    local service_name="$1"

    while IFS= read -r line; do
        # Skip all HTTP 200 status requests (successful requests - usually not needed for debugging)
        if echo "$line" | grep -q '"status":200'; then
            continue
        fi

        # Skip HTTP 401 status requests (authentication failures - often just expired sessions)
        if echo "$line" | grep -q '"status":401'; then
            continue
        fi

        # Skip repetitive x-casaos errors (show only once at startup)
        if echo "$line" | grep -q "extension \`x-casaos\` not found"; then
            local error_marker="/tmp/xcasaos-error-${service_name}"
            if [ ! -f "$error_marker" ]; then
                echo "$line"
                touch "$error_marker"
            fi
            continue
        fi

        # Skip repetitive NVIDIA GPU errors (show only once at startup)
        if echo "$line" | grep -q "NvidiaGPUInfoList error.*nvidia-smi.*executable file not found"; then
            local gpu_marker="/tmp/nvidia-error-${service_name}"
            if [ ! -f "$gpu_marker" ]; then
                echo "$line"
                touch "$gpu_marker"
            fi
            continue
        fi

        # Skip Chinese ping messages
        if echo "$line" | grep -q "消息来了"; then
            continue
        fi

        # Skip any line that's just a user agent string (these are often incomplete log lines)
        if echo "$line" | grep -q "Mozilla/5.0.*Chrome.*Safari" && ! echo "$line" | grep -q '"time":'; then
            continue
        fi

        # Skip health check endpoints
        if echo "$line" | grep -q '"uri":".*health"'; then
            continue
        fi

        # Skip ping/heartbeat endpoints
        if echo "$line" | grep -q '"uri":".*ping"'; then
            continue
        fi

        # Skip WebSocket upgrade requests (unless they're errors)
        if echo "$line" | grep -q '"uri":".*websocket"' && echo "$line" | grep -q '"status":101'; then
            continue
        fi

        # Show everything else (actual errors, warnings, important info)
        echo "$line"
    done
}

exec gosu "$PUID:$PGID" /usr/local/bin/casaos-main 2>&1 | filter_logs