# s6-overlay/s6-rc.d/register-ui-events/up
#!/command/with-contenv bash

chown -R "$PUID:$PGID" /usr/local/bin/register-ui-events.sh
exec gosu "$PUID:$PGID" /usr/local/bin/register-ui-events.sh