# s6-overlay/s6-rc.d/register-ui-events/up
/etc/s6-overlay/scripts/register-ui-events.sh